# beholder-otg/beholder_docker.service
[Unit]
Description=Beholder
Requires=docker.service
After=docker.service beholder_otg.service beholder_boot.service

[Service]
Type=oneshot
RemainAfterExit=true
WorkingDirectory=/home/beholder/beholder/beholder-otg/
ExecStart=beholder_docker.sh up
ExecStop=beholder_docker.sh down
Restart=on-failure

[Install]
WantedBy=multi-user.target